# src/conf/fty-email.service.in
[Unit]
Description=@PROJECT_NAME@ service
Requires=network.target malamute.service
After=network.target malamute.service
#WAS#Before=fty-asset.service
After=fty-asset.service
Requires=fty-asset.service
BindsTo=fty-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@CMAKE_INSTALL_FULL_DATAROOTDIR@/bios/etc/default/bios
EnvironmentFile=-@CMAKE_INSTALL_FULL_DATAROOTDIR@/bios/etc/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_FULL_DATAROOTDIR@/fty/etc/default/fty
EnvironmentFile=-@CMAKE_INSTALL_FULL_DATAROOTDIR@/fty/etc/default/fty__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/bios
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/fty
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/default/fty__%n.conf
Environment="prefix=@CMAKE_INSTALL_FULL_PREFIX@"
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/@PROJECT_NAME@ --config @CMAKE_INSTALL_FULL_SYSCONFDIR@/@PROJECT_NAME@/@PROJECT_NAME@.cfg

[Install]
WantedBy=bios.target
